# target/linux/armsr/base-files/etc/uci-defaults/kmods-unload-armsr
#!/bin/sh

[ -f /etc/config/kmods ] || exit 0

uci -q batch <<-EOF >/dev/null

	delete kmods.mwifiex
	set kmods.mwifiex=auto_unload
	set kmods.mwifiex.function=mwifiex_pcie
	set kmods.mwifiex.support=mwifiex
	delete kmods.ath10k
	set kmods.ath10k=auto_unload
	set kmods.ath10k.function=ath10k_pci
	set kmods.ath10k.support=ath10k_core

	delete kmods.qed
	set kmods.qed=auto_unload
	set kmods.qed.function="qede qedf"
	set kmods.qed.support=qed

	delete kmods.armsr_others
	set kmods.armsr_others=auto_unload
	set kmods.armsr_others.function="atlantic bnx2x i40e ixgbe mlx5_core mlx4_core mpt3sas sfc sfc_falcon"

	commit kmods
EOF

exit 0
